# config/config-files.list
/var/lib/alsa/asound.state
/etc/direwolf/direwolf.conf
/etc/asound.conf
/etc/polaric-webapp
/etc/polaric-aprsd
/etc/apache2/sites-available/logio.conf
/etc/init.d/direwolf
/etc/init.d/log.io-harvester
/etc/init.d/log.io-server
/etc/init.d/polaric-aprsd
/etc/init.d/polaric-webapp
/etc/logrotate.d/direwolf
/etc/logrorate.d/polaric-aprsd
/etc/rc2.d/S001direwolf
/etc/rc2.d/S01polaric-aprsd
/etc/rc2.d/S01polaric-webapp
/etc/rsyslog.d/direwolf.conf


